# server-config-files/supervisor-conf.d-mytornado.conf
[program:smart-mirror-8888]
command=sudo python /home/pi/dev/smart-mirror/python/app.py --port=8888 
stderr_logfile = /var/log/tornado-stderr.log
stdout_logfile = /var/log/tornado-stdout.log
